# infra/templates/nebula.service.j2
[Unit]
Description=Nebula Network
Wants=basic.target
After=basic.target network.target

[Service]
ExecStart=/usr/bin/nebula -config /etc/nebula/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
